# .crontab
MAILTO=""
*/5 * * * * offlineimap-notify
*/5 * * * * updatedb -l 0 -U $HOME -o $HOME/.updatedb
